# recipes-extended/libyocto/libyocto-config/systemd/chronyd.service.d/90-fareco.conf
[Unit]
ConditionPathExists=/etc/fareco/services/ntp
