# recipes-extended/libyocto/libyocto-config/systemd/wpa_supplicant-wired@.service.d/90-fareco.conf
[Unit]
ConditionPathExists=/etc/fareco/services/wpa-wired

[Service]
ExecStart=
ExecStart=/usr/sbin/wpa_supplicant -c/etc/fareco/wpa_supplicant/wpa_supplicant-wired-%I.conf -Dwired -i%I

