# recipes-extended/libyocto/libyocto-config/systemd/systemd-journal-upload.service.d/90-fareco.conf
[Unit]
ConditionPathExists=/etc/fareco/services/journal-upload
